# config_docker/supervisord.conf
[program:disquaire-gunicorn]
command = /home/app/env/bin/gunicorn gestopers.wsgi:application
user = root
directory = /home/app
autostart = true
autorestart = true
environment = ENV="test"